# packages/aoe2-telegram-bot/aoe2_telegram_bot-0.2.2.tar.gz/aoe2_telegram_bot-0.2.2/src/aoe2_telegram_bot/distro/aoe2-telegram-bot.service
[Unit]
Description=Age Of Empires 2 Telegram Bot
After=network.target

[Service]
Type=simple
User=%i
Group=%i

EnvironmentFile=%h/.config/aoe2-telegram-bot/env
ExecStart=/usr/bin/python3 -m aoe2_telegram_bot

Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
