# packages/tinybeans-sync/tinybeans_sync-0.2.tar.gz/tinybeans_sync-0.2/systemd/tinybeans-sync.timer
[Unit]
Description=Sync Tinybeans photos timer
Requires=tinybeans-sync.service

[Timer]
OnBootSec=15min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
